# clojure/.config/systemd/user/doctor-topbar.service
[Unit]
Description=doctor topbar app

[Service]
Type=simple
ExecStart=bb --config /home/russ/russmatney/clawe/bb.edn topbar
Environment="DISPLAY=:0"
Restart=on-failure

[Install]
WantedBy=default.target
