# clojure/.config/systemd/user/doctor-be.service
[Unit]
Description=doctor server

[Service]
Type=simple
Environment="PATH=/home/russ/.local/bin:/home/russ/.nix-profile/bin:/bin:/usr/bin:/usr/local/bin:/usr/local/sbin"
Environment="DISPLAY=:0"
ExecStart=bb --config /home/russ/russmatney/clawe/bb.edn doctor-be
Restart=on-failure

[Install]
WantedBy=default.target
