# clojure/.config/systemd/user/ralphie-nrepl.service
[Unit]
Description=ralphie babashka nrepl

[Service]
Type=simple
ExecStart=/home/russ/.local/bin/ralphie-nrepl

[Install]
WantedBy=graphical.target
